# templates/default/heartbeat.conf.erb
pre-start script

end script

start on (local-filesystems and net-device-up and runlevel [2345])
stop on runlevel [!2345]

respawn
respawn limit 2 10   # stop if there are 2 respawns in 10 seconds

setuid root
setgid root

env CONTROLLER_IP=<%= @controller_ip %>

script
  exec ruby /opt/deis/runtime/heartbeat.rb &> /var/log/deis/heartbeat.log
end script

