# templates/default/rendevous.conf.erb
pre-start script

end script

start on (local-filesystems and net-device-up and runlevel [2345])
stop on runlevel [!2345]

respawn
respawn limit 2 10   # stop if there are 2 respawns in 10 seconds

setuid deis
setgid docker

script
  cd /opt/deis/controller
  . venv/bin/activate
  exec python bin/rendevous.tpy >> /var/log/deis/rendevous.log 2>&1
end script

